# templates/druid-coordinator.service.erb
[Unit]
Description=Druid Coordinator Node

[Service]
Type=simple
StandardOutput=syslog
StandardError=syslog
SyslogFacility=<%= scope.lookupvar("::#{@module_name}::syslog_facility") %>
WorkingDirectory=<%= scope.lookupvar("::#{@module_name}::install_dir") %>/druid/
ExecStart=/usr/bin/java <%= @jvm_opts.join(' ') %> -classpath .:/etc/druid/coordinator:/<%= scope.lookupvar("::#{@module_name}::extra_classpaths") %>:<%= scope.lookupvar("::#{@module_name}::install_dir") %>/druid/lib/* org.apache.druid.cli.Main server coordinator
SuccessExitStatus=130 143
Restart=on-failure

[Install]
WantedBy=multi-user.target
